# _misc/arch/etc/local.d/local.start
# /etc/conf.d/local.start
# This is a good place to load any misc programs
# on startup (use &>/dev/null to hide output)

#. /etc/conf.d/ntpclient

#echo " * Starting firewall ..."
#/etc/conf.d/firewall.sh start

[ -e /mnt/vmShared ] || mkdir -p /mnt/vmShared
fsck.ext4 /dev/sdb1
mount /dev/sdb1 /mnt/hd/
/usr/bin/vmhgfs-fuse -o allow_other -o auto_unmount .host:/vmShared /mnt/vmShared

/sbin/shmtool start

#sudo -u tarjanm screen -S mldonkey -d -m nice mlnet
