# usr/lib/systemd/system/youtube-dl-upd.timer
[Unit]
Description=Daily youtube-dl update check

[Timer]
OnCalendar=daily
AccuracySec=12h
Persistent=true

[Install]
WantedBy=timers.target
