# _misc/arch/etc/profile.d/proxy.sh
export http_proxy="http://proxy.server.com:8080"
export https_proxy="https://proxy.server.com:8443"
export ftp_proxy="http://proxy.server.com:1080"
export rsync_proxy=$http_proxy
export no_proxy="localhost,127.0.0.1,archvm"
#,.movidius.com
