# etc/local.d/local.stop
# /etc/conf.d/local.stop

# This is a good place to unload any misc.
# programs you started above.
# For example, if you are using OSS and have
# "/usr/local/bin/soundon" above, put
# "/usr/local/bin/soundoff" here.

. /etc/conf.d/brightness

echo 0 > $BRTNS_DEV
nc 127.0.0.1 4000 -ve "echo kill" | logger
mv -f /home/m/.mozilla/firefox/ramdisk.tgz /home/m/.mozilla/firefox/ramdisk.tgz~
tar c -C /dev/shm/firefox `ls -1 /dev/shm/firefox` | pigz > /home/m/.mozilla/firefox/ramdisk.tgz
